# etc/services.d/fhem/run
#!/bin/bash

#umask 022

cd /app/fhem

 exec \
 	s6-setuidgid abc perl /app/fhem/fhem.pl /app/fhem/fhem.cfg
